# roles/expoter_install/files/node-exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart=/usr/bin/nohup /usr/sbin/node_exporter

[Install]
WantedBy=multi-user.target
